# tools/systemd/bazel-matrix-project-loop.service
[Unit]
Description=Run one Bazel Matrix project addition or triage step
Documentation=https://github.com/dan-stowell/bazel-matrix

[Service]
Type=oneshot
WorkingDirectory=/home/exedev/bazel-matrix
ExecStart=/home/exedev/bazel-matrix/tools/run-one-project-loop.sh
RuntimeMaxSec=6h 30min
KillSignal=SIGINT
TimeoutStopSec=2min
